# .config/systemd/user/randomwallpaper.timer
[Unit]
Description=Change the wallpaper every 20 mins

[Timer]
OnBootSec=1min
OnUnitActiveSec=20min

[Install]
WantedBy=timers.target

