# .config/systemd/user/randomwallpaper.timer
[Unit]
Description=Change the wallpaper ever half hour

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target

